# scripts/com.llmxfactors.scout.discover.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
    Template — install_launchd.sh substitutes the placeholders at install
    time so the generated plist matches the actual user, project location,
    and Homebrew prefix on the target Mac. Do NOT install this file directly.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.llmxfactors.scout.discover</string>

    <key>ProgramArguments</key>
    <array>
        <string>@@UV_BIN@@</string>
        <string>run</string>
        <string>scout</string>
        <string>discover</string>
        <string>--limit</string>
        <string>50</string>
    </array>

    <key>WorkingDirectory</key>
    <string>@@PROJECT_ROOT@@</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>@@PATH@@</string>
        <key>HOME</key>
        <string>@@HOME@@</string>
    </dict>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key><integer>9</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>15</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>21</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>AbandonProcessGroup</key>
    <false/>

    <key>StandardOutPath</key>
    <string>@@PROJECT_ROOT@@/logs/cron.log</string>
    <key>StandardErrorPath</key>
    <string>@@PROJECT_ROOT@@/logs/cron.log</string>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
